# facerec_webpage.timer
[Unit]
Description=Timer for facerec webpage

[Timer]
OnUnitInactiveSec=10
Unit=facerec_webpage.service

[Install]
WantedBy=timers.target
